# usb-cleaner-box.service
[Unit]
Description=USB Cleaner Box - Security Scanner and Cleaner
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/pi/UCB
ExecStart=/usr/bin/python3 /home/pi/UCB/usb_cleaner_box.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
